# test/exec-privatedevices-no.service
[Unit]
Description=Test for PrivateDev=no

[Service]
ExecStart=/bin/sh -c 'exit $(test -c /dev/kmsg)'
Type=oneshot
PrivateDevices=no
