# cmd/package/unitfiles/auxiliary-manager.target
[Unit]
Description=Auxiliary Service Management
Requires=basic.target
Conflicts=rescue.service rescue.target
After=basic.target rescue.service rescue.target
AllowIsolate=yes
